# systemd/nightly-dream.service
[Unit]
Description=codeREM nightly memory consolidation
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/nightly-dream
Environment=HOME=%h
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
